# scripts/crontabs/root_crontab
0 15 * * 1,3,5 /usr/bin/pacman-mirrors -f >> /home/alux/workbench/System/Logging\ and\ Rollback/paclog.log

